# home-assistant/integrations/cron/home-assistant.cron
# Home Assistant Backups


* * * * * /bin/bash /config/integrations/shell_scripts/trigger-ha-api.sh "start_github_config_share"

