# home-assistant/integrations/cron/pfsense.cron
# PFSense Backups
* * * * * bash /config/integrations/shell_scripts/trigger-ha-api.sh "start_pfsense_backup"