# home-assistant/integrations/cron/vault1.cron
#### Power On #############
# 6:55PM Every 3 days
55 18 * * * bash /config/integrations/shell_scripts/trigger-ha-api.sh "start_vault1"

# 10:55 AM On 20th of every month
55 10 20 * * bash /config/integrations/shell_scripts/trigger-ha-api.sh "start_vault1"

#### Power Off #############
# 03:30 PM On 20th of every month
30 15 20 * * bash /config/integrations/shell_scripts/trigger-ha-api.sh "stop_vault1"

# 10 PM Every Day
0 21 * * * bash /config/integrations/shell_scripts/trigger-ha-api.sh "stop_vault1"